# ramdisk/mnt/etc/inittab
::respawn:/bin/sh
::respawn:/sbin/syslogd -n
::sysinit:/etc/rc
#::bootwait:/bin/cat /etc/motd
